# server/vps/mtproto-web-tls-reload.service
[Unit]
Description=Restart MTProto Admin Web TLS when certificate files change

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart mtproto-web-tls.service

