# server/vps/mtproto-traffic-sync.service
[Unit]
Description=Sync MTProto traffic iptables rules
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mtproto-traffic sync

[Install]
WantedBy=multi-user.target
